# crontab
# m h  dom mon dow   command
0 0,4,8,12,16,20 * * * root  . /tmp/env.sh; python3 /seneye/seneye-offline-check.py > /proc/1/fd/1 2>/proc/1/fd/2
